# systemd/bonescript.socket
[Socket]
ListenStream=80

[Install]
WantedBy=sockets.target
