# etc/profile.d/node.sh
#!/bin/sh
export NODE_PATH=/usr/lib/node_modules
